# services/kairos-turret-test.service
[Unit]
Description=KaiROS Servo Daemon compatibility unit (deprecated)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/kairos
ExecStart=/opt/kairos/apps/bot/servo_daemon
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
